# xenwooting/contrib/systemd/xenwooting.service
[Unit]
Description=XenWooting (Wooting -> WTN -> MTS-ESP + MIDI + RGB)
After=local-fs.target
Wants=local-fs.target

# If MODEP is present, start before the audio host so clients see the master at boot.
Before=modep-mod-ui.service
Before=modep-mod-host.service

# Avoid running multiple masters / conflicting consumers.
Conflicts=tun.service
Conflicts=wooting-analog-midi-headless.service

[Service]
Type=simple
User=patch
Group=patch
WorkingDirectory=/home/patch/wooting-xen/xenwooting
Environment=RUST_LOG=info

# Keep the web UI control-plane on tmpfs (/run) so filesystem stalls cannot block the realtime loop.
Environment=XENWTN_RUNTIME_DIR=/run/xenwooting
Environment=XENWTN_LIVE_STATE_PATH=/run/xenwooting/live.json
Environment=XENWTN_PREVIEW_ENABLED_PATH=/run/xenwooting/preview.enabled
Environment=XENWTN_PREVIEW_WTN_PATH=/run/xenwooting/preview.wtn
Environment=XENWTN_TRAINER_MODE_PATH=/run/xenwooting/trainer.mode
Environment=XENWTN_GUIDE_PATH=/run/xenwooting/guide.json
Environment=XENWTN_HIGHLIGHT_PATH=/run/xenwooting/highlight.txt

ExecStart=/usr/local/bin/xenwooting
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
